# rc.local
#!/bin/sh -e
sh /home/vinicius/start.sh
